# launchd/com.openbrain.mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Template: run scripts/install-services.sh to generate with your actual paths -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openbrain.mcp</string>
    <key>ProgramArguments</key>
    <array>
        <string>__BUN__</string>
        <string>run</string>
        <string>src/index.ts</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__REPO__</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__HOME__/.bun/bin:__HOME__/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>GOG_KEYRING_PASSWORD</key>
        <string>__GOG_KEYRING_PASSWORD__</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__REPO__/logs/mcp.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO__/logs/mcp.log</string>
</dict>
</plist>
