# launchd/com.openbrain.poll.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Template: run scripts/install-services.sh to generate with your actual paths -->
<!-- Headless source-sync trigger: POSTs /api/sources/poll-due on a timer so
     ingestion does not depend on the OpenBrain.app menubar app being open.
     Safe alongside the app: poll-due has an in-flight guard + concurrency cap. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openbrain.poll</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/curl</string>
        <string>-s</string>
        <string>-X</string>
        <string>POST</string>
        <string>http://localhost:6277/api/sources/poll-due</string>
        <string>-H</string>
        <string>Content-Type: application/json</string>
    </array>
    <!-- Fire 3x/day: morning, midday, early evening. launchd is the SOLE
         heartbeat for /api/sources/poll-due (the tray app no longer runs its
         own timer), so these three times are the only moments sources sync —
         and thus the only moments Firecrawl is hit. Missed fires (Mac asleep)
         run once on wake. RunAtLoad intentionally omitted to avoid an extra
         off-window scrape at login. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>7</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>12</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>18</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>
    <key>StandardOutPath</key>
    <string>__REPO__/logs/poll.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO__/logs/poll.log</string>
</dict>
</plist>
